# ollama/etc/s6-overlay/s6-rc.d/serve/run
#!/command/with-contenv bash

ollama serve
